# root/var/lib/dpkg/info/rpi-eeprom.list
/.
/etc
/etc/default
/etc/default/rpi-eeprom-update
/lib
/lib/firmware
/lib/firmware/raspberrypi
/lib/firmware/raspberrypi/bootloader-2711
/lib/firmware/raspberrypi/bootloader-2711/default
/lib/firmware/raspberrypi/bootloader-2711/default/pieeprom-2020-04-16.bin
/lib/firmware/raspberrypi/bootloader-2711/default/pieeprom-2020-09-03.bin
/lib/firmware/raspberrypi/bootloader-2711/default/pieeprom-2021-03-18.bin
/lib/firmware/raspberrypi/bootloader-2711/default/pieeprom-2021-04-29.bin
/lib/firmware/raspberrypi/bootloader-2711/default/pieeprom-2022-01-25.bin
/lib/firmware/raspberrypi/bootloader-2711/default/pieeprom-2022-04-26.bin
/lib/firmware/raspberrypi/bootloader-2711/default/pieeprom-2022-11-25.bin
/lib/firmware/raspberrypi/bootloader-2711/default/pieeprom-2022-12-07.bin
/lib/firmware/raspberrypi/bootloader-2711/default/pieeprom-2023-01-11.bin
/lib/firmware/raspberrypi/bootloader-2711/default/recovery.bin
/lib/firmware/raspberrypi/bootloader-2711/default/vl805-000137ad.bin
/lib/firmware/raspberrypi/bootloader-2711/default/vl805-000138a1.bin
/lib/firmware/raspberrypi/bootloader-2711/default/vl805-000138c0.bin
/lib/firmware/raspberrypi/bootloader-2711/latest
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2020-04-16.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2020-06-15.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2020-07-16.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2020-07-31.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2020-09-03.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2020-12-11.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2021-01-11.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2021-01-16.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2021-02-16.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2021-03-18.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2021-04-29.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2021-07-06.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2021-11-22.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2021-12-02.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2022-01-25.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2022-02-08.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2022-03-10.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2022-04-26.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2022-07-19.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2022-07-22.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2022-07-26.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2022-08-02.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2022-09-02.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2022-10-18.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2022-11-04.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2022-11-25.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2022-12-07.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2023-01-11.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/pieeprom-2023-05-11.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/recovery.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/vl805-000137ad.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/vl805-000138a1.bin
/lib/firmware/raspberrypi/bootloader-2711/latest/vl805-000138c0.bin
/lib/firmware/raspberrypi/bootloader-2711/release-notes.md
/lib/firmware/raspberrypi/bootloader-2712
/lib/firmware/raspberrypi/bootloader-2712/default
/lib/firmware/raspberrypi/bootloader-2712/default/pieeprom-2023-09-13.bin
/lib/firmware/raspberrypi/bootloader-2712/default/pieeprom-2023-09-21.bin
/lib/firmware/raspberrypi/bootloader-2712/default/pieeprom-2023-09-28.bin
/lib/firmware/raspberrypi/bootloader-2712/default/pieeprom-2023-10-18.bin
/lib/firmware/raspberrypi/bootloader-2712/default/pieeprom-2023-10-30.bin
/lib/firmware/raspberrypi/bootloader-2712/default/pieeprom-2023-11-20.bin
/lib/firmware/raspberrypi/bootloader-2712/default/recovery.bin
/lib/firmware/raspberrypi/bootloader-2712/latest
/lib/firmware/raspberrypi/bootloader-2712/latest/pieeprom-2023-09-13.bin
/lib/firmware/raspberrypi/bootloader-2712/latest/pieeprom-2023-09-21.bin
/lib/firmware/raspberrypi/bootloader-2712/latest/pieeprom-2023-09-28.bin
/lib/firmware/raspberrypi/bootloader-2712/latest/pieeprom-2023-10-18.bin
/lib/firmware/raspberrypi/bootloader-2712/latest/pieeprom-2023-10-30.bin
/lib/firmware/raspberrypi/bootloader-2712/latest/pieeprom-2023-11-20.bin
/lib/firmware/raspberrypi/bootloader-2712/latest/recovery.bin
/lib/firmware/raspberrypi/bootloader-2712/release-notes.md
/lib/systemd
/lib/systemd/system
/lib/systemd/system/rpi-eeprom-update.service
/usr
/usr/bin
/usr/bin/rpi-eeprom-config
/usr/bin/rpi-eeprom-digest
/usr/bin/rpi-eeprom-update
/usr/share
/usr/share/doc
/usr/share/doc/rpi-eeprom
/usr/share/doc/rpi-eeprom/LICENCE.bootloader
/usr/share/doc/rpi-eeprom/changelog.Debian.gz
/usr/share/doc/rpi-eeprom/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/rpi-eeprom-config.1.gz
/usr/share/man/man1/rpi-eeprom-update.1.gz
/var
/var/lib
/var/lib/raspberrypi
/var/lib/raspberrypi/bootloader
/var/lib/raspberrypi/bootloader/backup
/lib/firmware/raspberrypi/bootloader-2711/beta
/lib/firmware/raspberrypi/bootloader-2711/critical
/lib/firmware/raspberrypi/bootloader-2711/stable
/lib/firmware/raspberrypi/bootloader-2712/beta
/lib/firmware/raspberrypi/bootloader-2712/critical
/lib/firmware/raspberrypi/bootloader-2712/stable
